# stage2-0/00-install-packages/files/etc/cron.d/checkwlan0
* * * * * root /usr/local/bin/check-wlan0.sh | /usr/bin/logger -t check-wlan0

